# contrib/systemd/anet-server.service
[Unit]
Description=ANet VPN Server
After=network.target
Wants=network.target
Before=network-online.target

[Service]
Type=notify
NotifyAccess=all
DynamicUser=yes
SupplementaryGroups=net-admin
NoNewPrivileges=yes

RootDirectory=/etc/anet
BindPaths=/dev/net/tun
BindReadOnlyPaths=/etc/anet/server.toml
ReadWritePaths=/var/log/anet
ReadWritePaths=/run/anet
TemporaryFileSystem=/tmp:/var/tmp
ProtectSystem=strict

RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
IPAccounting=yes

AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_IPC_LOCK
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_IPC_LOCK

LogNamespace=anet
StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=info

WorkingDirectory=/etc/anet
ExecStartPre=/sbin/modprobe tun
ExecStart=/usr/bin/anet-server -c /etc/anet/server.toml
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
TimeoutStartSec=30s

[Install]
WantedBy=multi-user.target
